# config/test.service
[Unit]
Description=test

[Service]
Type=simple
Restart=always
RestartSec=5s
ExecStart=/usr/local/go/bin/test

[Install]
WantedBy=multi-user.target
